# packages/systemd/desktop-security-center.service
[Unit]
Description=Desktop security center backend
Requires=desktop-security-center.socket

[Service]
Type=simple
ExecStart=/usr/bin/desktop-security-center

[Install]
WantedBy=multi-user.target
