# packages/systemd/desktop-security-center.socket
[Unit]
Description=Socket for desktop-security-center backend/frontend communication

[Socket]
ListenStream=%t/desktop-security-center.socket

[Install]
WantedBy=sockets.target
