# packages/systemd/server.service
[Unit]
Description=Whatever
Requires=network.target
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/seccenter

[Install]
WantedBy=multi-user.target

